# 部署工具/systemd/crypto-market-microstructure.service
[Unit]
Description=Crypto Market Microstructure Feature Collector
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/crypto-auto-trader
EnvironmentFile=-/etc/crypto-auto-trader/trading.env
ExecStart=/opt/crypto-auto-trader/.venv/bin/python -B /opt/crypto-auto-trader/market_microstructure_service.py --root /opt/crypto-auto-trader --interval 30 --top-limit 100 --batch-size 15
Restart=always
RestartSec=5
Environment=PYTHONIOENCODING=utf-8
Environment=OKX_MARKET_DATA_ENABLED=1
Environment=BYBIT_MARKET_DATA_ENABLED=1
Environment=OKX_MARKET_DATA_MAX_PER_MIN=60
Environment=BYBIT_MARKET_DATA_MAX_PER_MIN=60
Environment=MICROSTRUCTURE_RETENTION_DAYS=14
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
