# 部署工具/systemd/crypto-account-state.service
[Unit]
Description=Crypto Central Account State Baseline Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/crypto-auto-trader
EnvironmentFile=/etc/crypto-auto-trader/trading.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/crypto-auto-trader/.venv/bin/python -u /opt/crypto-auto-trader/account_state_service.py --bootstrap-empty

[Install]
WantedBy=multi-user.target
