# 部署工具/systemd/crypto-market-data-cache.service
[Unit]
Description=Crypto Market Data Cache Service
After=network-online.target crypto-binance-api-queue.service
Wants=network-online.target crypto-binance-api-queue.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/crypto-auto-trader
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/crypto-auto-trader/.venv/bin/python -u /opt/crypto-auto-trader/market_data_service.py --root /opt/crypto-auto-trader --interval 300 --top-limit 180
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
