# 部署工具/systemd/crypto-data-maintenance.timer
[Unit]
Description=Hourly Crypto Data Retention Maintenance Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
