# 部署工具/systemd/crypto-binance-user-stream-v14.service
[Unit]
Description=AutoTrading Binance user-data-stream service C/v14 (construction foundation)
After=network-online.target crypto-binance-api-queue.service crypto-account-state.service
Wants=network-online.target crypto-binance-api-queue.service crypto-account-state.service

[Service]
Type=simple
WorkingDirectory=/opt/crypto-auto-trader
EnvironmentFile=-/etc/crypto-auto-trader/trading.env
ExecStart=/opt/crypto-auto-trader/.venv/bin/python -u /opt/crypto-auto-trader/binance_user_stream_service.py --strategy C/v14 --account C --apply-state
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
